# supervisord.conf
[supervisord]
nodaemon=true

[program:myapp]
command=%(ENV_MYAPP_COMMAND)s
autostart=true
autorestart=true
startretries=5


stderr_logfile=/var/log/supervisord.err.log
stdout_logfile=/var/log/supervisord.out.log
